# system/overlays/config/etc/systemd/system/systemd-firstboot.service.d/override.conf
# Prompt everything
[Service]
ExecStart=
ExecStart=systemd-firstboot --prompt --hostname=openwrt_one --root-password=toor --root-shell=/bin/bash
